# archiso.config/root-image/etc/systemd/system/CoinDrone.service
[Unit]
Description=CoinDrone Process

[Service]
Type=forking
ExecStart=/etc/systemd/scripts/CoinDrone

[Install]
WantedBy=multi-user.target
